# trend-bot.service
[Unit]
Description=BTC Crypto Trend Following Trading Bot
After=network.target

[Service]
Type=simple
User=cryptbot
WorkingDirectory=/home/cryptbot/cryptbot
ExecStart=/home/cryptbot/cryptbot/venv/bin/python /home/cryptbot/cryptbot/main.py
Restart=always
RestartSec=60
StandardOutput=append:/home/cryptbot/cryptbot/logs/systemd.log
StandardError=append:/home/cryptbot/cryptbot/logs/systemd.log
Environment="PATH=/home/cryptbot/cryptbot/venv/bin:/usr/bin:/bin"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
